# launchagent.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>__LABEL__</string>
  <key>ProgramArguments</key>
  <array>
    <string>__PYTHON__</string>
    <string>__APP_DIR__/claude_token_ring.py</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>
  <key>WatchPaths</key>
  <array>
    <string>__HOME__/Library/Application Support/Claude/Cache/Cache_Data</string>
  </array>
  <key>StandardOutPath</key>
  <string>/tmp/claude-token-ring__NS_LABEL__.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/claude-token-ring__NS_LABEL__.log</string>
</dict>
</plist>
